# microscope_control.service
[Unit]
Description=Microscope UDP controller service
After=network-online.target

[Service]
ExecStart=$MICROSCOPE_CONTROLLER_HOME/udp_controller.py
User=pi

[Install]
WantedBy=multi-user.target
